# roles/vm_drop_caches/files/drop_caches.service
[Unit]
Description=Drop virtual memory caches

[Service]
ExecStart=/bin/sh -c 'sync; echo 3 > /proc/sys/vm/drop_caches'
