# roles/docker2nftables/templates/docker2nftables.service
[Unit]
Description=docker2nftables forward rules

[Service]
ExecStart=/usr/local/bin/docker2nftables
